# cgminer/root-files/etc/rc.local

# Load hasher bitstream
[ -e /dev/xdevcfg ] && cat /root/hasher.bit > /dev/xdevcfg

# Close ipv6
sleep 3
sysctl -p

exit 0

